# systemd/audio-watcher.service
[Unit]
Description=Watch for USB audio device changes and refresh PipeWire
After=pipewire.service wireplumber.service

[Service]
Type=simple
ExecStart=%h/dotfiles/hypr/scripts/audio_watcher.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
